# virtualization-systems/cw3/ex2/service/tvsctld.socket
[Unit]
Description=Tvsctld Activation Socket.

[Socket]
ListenStream=/run/isel/tvsctld/request

[Install]
WantedBy=sockets.target
